# deploy/mc-loki/conductor-mc.service
[Unit]
Description=Conductor Master Controller
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=nthmost
WorkingDirectory=/opt/conductor-mc
ExecStart=/usr/bin/java \
    -jar /opt/conductor-mc/conductor-server-boot.jar \
    --spring.config.location=/opt/conductor-mc/application.properties
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=conductor-mc

[Install]
WantedBy=multi-user.target
